# backend/dtcc-table.service
[Unit]
Description=DTCC Table Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/dtcc-table/backend
Environment="PATH=/home/ubuntu/dtcc-table/backend/venv/bin"
EnvironmentFile=/home/ubuntu/dtcc-table/backend/.env
ExecStart=/home/ubuntu/dtcc-table/backend/venv/bin/python app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target